# deploy/systemd/skill-hub-sync.service
[Unit]
Description=Sync Skill Hub release

[Service]
Type=oneshot
WorkingDirectory=/opt/skill-hub
ExecStart=/bin/sh -lc 'echo "replace this with artifact download and current symlink switch"'
